# src/vm/data/traefik.service
[Unit]
Description=Traefik Service
After=network.target

[Service]
ExecStart=/etc/traefik/launch.sh

[Install]
WantedBy=multi-user.target
